# crm.sh
#!/bin/bash
#chkconfig: 2345 80 90
#description: crm
cd /opt/crm-dit-linux-x64
./CRM.UI &
